# CONTAINER-DEBIAN-PKG/etc/systemd/system/automation-interface-api-f5-container.service
[Unit]
Requires=network-online.target
After=automation-interface-dns-container.service

[Service]
Type=oneshot
ExecStart=/usr/bin/api-f5-container.sh start
ExecStop=/usr/bin/api-f5-container.sh stop
KillMode=process
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
